# packaging/bitprotector.service
[Unit]
Description=BitProtector - Distributed File Mirror and Integrity Protection
After=network.target
Documentation=https://github.com/bitprotector/bitprotector

[Service]
Type=simple
User=bitprotector
Group=bitprotector
ExecStart=/usr/bin/bitprotector \
    --db /var/lib/bitprotector/bitprotector.db \
    serve \
    --host 127.0.0.1 \
    --port 8443 \
    --tls-cert /etc/bitprotector/tls/cert.pem \
    --tls-key /etc/bitprotector/tls/key.pem
Restart=on-failure
RestartSec=5s

# Security hardening
NoNewPrivileges=true
# "full" keeps /usr, /boot, /efi, and /etc read-only while leaving data
# directories (/mnt, /media, /home, /run/media, etc.) writable — necessary
# because drive paths and virtual-path symlinks can be at any absolute path.
ProtectSystem=full
ReadWritePaths=/var/lib/bitprotector /var/log/bitprotector
# Allows the service to write to drive directories regardless of their
# ownership (e.g. root-owned ext4 mount roots), so users do not need to
# manually chown mount points after install.
AmbientCapabilities=CAP_DAC_OVERRIDE
CapabilityBoundingSet=CAP_DAC_OVERRIDE
PrivateTmp=true

[Install]
WantedBy=multi-user.target
